# infra/systemd/hwax-stack.service
# HWAX full stack — one user service that brings up the whole federation via the
# orchestrator (portal + sub-services + chat stack), updating code first.
# Installed to ~/.config/systemd/user/ by infra/scripts/install-systemd.sh.
#
# Why one unit (not per-service): the orchestrator (services.py) already knows tier
# order, health, auto-discovery, and SSH for remote hosts. systemd just needs to call
# `up --update` on boot and `down` on stop. Rootless: runs as the user; linger makes it
# survive logout / start on boot. /run/user/<uid> is required by rootless apptainer.

[Unit]
Description=HWAX full stack (portal + sub-services + chat) via orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__PORTAL__
# up --update is IDEMPOTENT: git-pull each service (vllm opts out), then start only what's
# down (already-up services are skipped), in tier order, health-polled. So start/restart
# never needlessly bounces heavy services (mxwp's 5 containers, etc.) — it just heals + updates.
# PATH is baked PER-BOX at install time: install-systemd.sh fills the PATH= line below with the
# installing user's $PATH (pnpm/node/git/etc., wherever that box keeps them). systemd's own env is
# minimal, so without this the heax `make build` (pnpm) fails with 127. No hardcoded path lives in
# this template — each machine bakes its own at install.
Environment="PATH=__PATH__"
ExecStart=__PY__ __PORTAL__/infra/scripts/services.py up --update
# No ExecStop: stopping the unit must NOT tear the whole federation down (that bounces heavy
# stacks and risks the vLLM EngineCore-zombie OOM). To really stop everything, run the
# orchestrator by hand:  ./infra/scripts/services.sh down
# vLLM model load + multi-service health can take a while on a cold boot.
TimeoutStartSec=900

[Install]
WantedBy=default.target
